# deploy/selfservice-janitor.timer
[Unit]
Description=Run the Self-Service Provisioning TTL sweep every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=1min

[Install]
WantedBy=timers.target
